# systemd/.config/systemd/user/kmonad@hhkb.path
[Unit]
Description=start kmonad for hkhb
StartLimitIntervalSec=400
StartLimitBurst=3

[Path]
PathExists=/dev/hhkb

[Install]
WantedBy=graphical-session.target%
