# systemd/.config/systemd/user/kmonad@internal_keyboard.service
[Unit]
Description=KMonad for internal keyboard
StartLimitIntervalSec=1 min
StartLimitBurst=3

[Service]
Type=simple
ExecStart=%h/.local/bin/kmonad %h/.config/kmonad/internal_keyboard.kbd
Restart=always
RestartSec=5
StandardOutput=null

[Install]
WantedBy=default.target
