# systemd/.config/systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -c 10 -w 30
